# deploy/systemd/nginx.service.d/xblp.conf
# systemd drop-in: order nginx after xblp-api so the TLS cert exists
# before nginx tries to load it. Without this, nginx fails on first boot
# because cert.pem is generated during xblp-api's startup lifespan.
[Unit]
After=xblp-api.service
Wants=xblp-api.service
