# deploy/systemd/xblp-api.service
[Unit]
Description=xboxlive-protect API daemon
Documentation=https://github.com/kstrube/xboxlive-protect/blob/main/docs/api-auth.md
After=network-online.target
Wants=network-online.target

# NOTE: xblp-capture.service does not exist yet (it will be added in Phase 2
# Stage 2 when the capture daemon's systemd integration is wired). When it
# does, consider adding: After=xblp-capture.service (but no Wants= — they are
# independent per the architectural decision; capture does not depend on the
# API being up and the API does not require capture to serve auth).

[Service]
Type=simple
User=xblp
Group=xblp
WorkingDirectory=/opt/xboxlive-protect
Environment=PYTHONPATH=/opt/xboxlive-protect/src PYTHONDONTWRITEBYTECODE=1 PYTHONUNBUFFERED=1
ExecStart=/opt/xboxlive-protect/venv/bin/python -m xblp_api
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xblp-api

# Creates /var/lib/xboxlive-protect owned by xblp:xblp before the service
# starts, so the API can open/create state.db without needing the install
# script to have run first. Subsumes ReadWritePaths for that directory.
StateDirectory=xboxlive-protect

# CAP_NET_ADMIN lets the daemon (running as xblp) apply the nftables ruleset
# at startup. Without it, apply_initial_ruleset() fails with EPERM on a fresh
# device that has no pre-existing inet xblp table.
AmbientCapabilities=CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_ADMIN
NoNewPrivileges=true

# Harden the filesystem view: /usr, /boot, /etc are read-only from the
# service's perspective. Using full rather than strict so Phase 3/5 can
# reach other /var paths without per-release access debugging.
ProtectSystem=full
# ProtectHome=true is intentionally not set in Stage 1 because the dev install
# layout symlinks code through /home/pi/Projects/xboxlive-protect into /opt/.
# Phase 5's image-build pipeline will deploy real code to /opt/xboxlive-protect
# and can reinstate ProtectHome=true at that point.
# Baseline security score without ProtectHome (Stage 1): 6.1 MEDIUM
# (systemd-analyze security xblp-api; ~5.8 expected after reinstatement)
PrivateTmp=true

[Install]
WantedBy=multi-user.target
